# dotfiles/.config/systemd/user/i3wm.target
[Unit]
Description=i3wm target
After=default.target
